# github-runner/github-runner.service
[Unit]
Description=GitHub Actions ephemeral runner loop
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/github-runner
EnvironmentFile=/home/pi/github-runner/.env
ExecStart=/home/pi/github-runner/runner-loop.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
